# deploy/nextdraw-api.service
[Unit]
Description=NextDraw Plotter API Server
After=network.target
Wants=network.target

[Service]
Type=simple
# Use the actual system username instead of hardcoded "pi"
User=%USER%
Group=%USER%
# Use %HOME% for proper home directory path
WorkingDirectory=%HOME%/plot-runner-agent
# Set environment variables for Python
Environment="PATH=%HOME%/plot-runner-agent/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=%HOME%/plot-runner-agent"
# Use a script that activates the virtual environment properly
ExecStart=%HOME%/plot-runner-agent/run_app.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%HOME%/plot-runner-agent
ProtectHome=yes

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextdraw-api

[Install]
WantedBy=multi-user.target
